# deploy/systemd/laravel-otp-whatsapp-queue@.service
[Unit]
Description=Laravel OTP WhatsApp Queue Worker %i
After=network.target

[Service]
Type=simple
User=litespeed
Group=litespeed
Restart=always
RestartSec=5
StartLimitInterval=0

# Worker dedicado para WhatsApp OTP con rate limiting
# Procesa máximo 5 mensajes por segundo según límite configurado
ExecStart=/usr/local/lsws/lsphp84/bin/php /var/www/webroot/ROOT/artisan queue:work redis \
    --queue=otp-whatsapp \
    --sleep=1 \
    --tries=3 \
    --max-jobs=5 \
    --max-time=1 \
    --memory=256 \
    --timeout=30

# Environment
Environment="QUEUE_CONNECTION=redis"
Environment="QUEUE_NAME=otp-whatsapp"

# Output
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target