# deploy/systemd/laravel-scheduler.service
[Unit]
Description=Laravel Scheduler
After=network.target redis.service

[Service]
Type=oneshot
User=nginx
Group=nginx
WorkingDirectory=/var/www/webroot/ROOT
ExecStart=/usr/bin/php /var/www/webroot/ROOT/artisan schedule:run
StandardOutput=append:/var/www/webroot/ROOT/storage/logs/scheduler.log
StandardError=append:/var/www/webroot/ROOT/storage/logs/scheduler.log